# install/kuma_agent.service
[Unit]
Description=Kuma Agent

[Service]
ExecStart=/usr/bin/kuma_agent
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target