# install/kuma_agent.service
[Unit]
Description=Kuma Agent

[Service]
ExecStart=/usr/local/bin/pc_manager
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target